# redis/supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/bin/bash start.sh
